# roles/daytradingbot/files/crypto_trader/ansible/roles/crypto-trader/templates/crypto-trader.service.j2
# /etc/systemd/system/crypto-trader.service — generated by Ansible

[Unit]
Description=Crypto Trader Webhook Server (Kraken Futures)
After=network-online.target cloudflared.service
Wants=network-online.target
# Start after cloudflared so the tunnel is up before we begin accepting signals

[Service]
User={{ trader_user }}
Group={{ trader_group }}
WorkingDirectory={{ trader_install_dir }}

ExecStart={{ trader_install_dir }}/venv/bin/uvicorn webhook.server:app \
    --host 127.0.0.1 \
    --port {{ trader_port }} \
    --workers 1 \
    --log-level info \
    --access-log

# ── Note: listen on 127.0.0.1, NOT 0.0.0.0 ──────────────────────────────────
# With Cloudflare Tunnel we don't expose any ports to the network.
# cloudflared connects to localhost:{{ trader_port }} internally.
# This means the webhook server is completely unreachable from outside
# unless the request comes through the Cloudflare tunnel.

Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-trader

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH={{ trader_install_dir }}

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
